# systemd/graeae-model-sync.service
[Unit]
Description=MNEMOS Model Registry — Daily Provider Sync
Documentation=https://github.com/your-org/mnemos
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=mnemos
Group=mnemos
WorkingDirectory=/opt/mnemos
ExecStart=/opt/mnemos/venv/bin/python3 scripts/sync_provider_models.py
Environment=DATABASE_URL=postgresql://mnemos:mnemos@localhost/mnemos
# Inherit API_KEYS_FILE from environment or rely on default ~/.api_keys_master.json
# Environment=API_KEYS_FILE=/etc/mnemos/api_keys.json

StandardOutput=journal
StandardError=journal
SyslogIdentifier=graeae-model-sync

# Allow partial failures (provider timeouts) without marking the unit failed
SuccessExitStatus=0 2
